# systemd/wifi_login.service
[Unit]
Description=Run wifi_login.sh every 120000 seconds

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi_login.sh


